# HPCC-Platform/build/initfiles/bash/etc/systemd/system/hpcc-environment-monitor.path
[Unit]
Description=monitor for /etc/HPCCSystems/environment.xml modification

[Path]
PathModified=/etc/HPCCSystems/environment.xml

[Install]
WantedBy=multi-user.target
