# HPCC-Platform/build/initfiles/bash/etc/systemd/system/thor-conf@.service
[Unit]
Description=generating %i configuration file and /var/lib/HPCCSystems/%i/%i.cfg
After=hpcc-conf@%i.service
Requires=hpcc-conf@%i.service

[Service]
Type=oneshot
User=hpcc
Group=hpcc
ExecStart=/opt/HPCCSystems/sbin/check-component-exists.sh %i
ExecStart=/bin/mkdir -p /var/lib/HPCCSystems/%i /var/log/HPCCSystems/%i /var/run/HPCCSystems /var/lock/HPCCSystems/%i
ExecStart=/bin/bash -c "/bin/sed -e 's/export //g' /var/lib/HPCCSystems/%i/setvars >> /var/lib/HPCCSystems/%i/%i.cfg"
